# Docker/cvesearch.cron
20 * * * *    python3 /cvesearch/sbin/db_updater.py -v >> /var/log/cvesearch.cron.log
