# Docker/cvesearch.cron
20 02 * * *    /cvesearch/Docker/updateDB.sh >> /var/log/cvesearch.cron.log 2>&1
